# fleet/v1/tests/fixtures/test.service
[Service]
ExecStart=/usr/bin/sleep 1d

